# docker-deps/supervisord.ini
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autorestart=true

[program:gunicorn]
directory=/var/www
command=/usr/bin/gunicorn --bind 0.0.0.0:8000 wsgi:app --workers 2 --keep-alive 2 --error-logfile /dev/stderr --access-logfile /dev/stdout --pid /var/run/gunicorn.pid
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0